# system_configuration/files/terminal_default_env/etc/profile.d/default_editor_pager.sh
## Set default text editor for console and terminals
export VISUAL=vim
export EDITOR=vim
export PAGER=less
